# root_config/etc/cron.weekly/fstrim
#!/bin/sh

# Uses trim on all filesystems listed in /etc/fstab that support it (SSD devices)
fstrim --all
