# root_config/etc/profile.d/environment.sh
# Set default editor for e.g `crontab -e`
export EDITOR=/usr/bin/nano
